# Packages/Admin/DeployStudio Admin.app/Contents/Frameworks/DSCore.framework/Versions/A/Resources/Tools/Finalize/com.deploystudio.finalizeScript.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.deploystudio.finalizeScript</string>
	<key>LimitLoadToSessionType</key>
	<string>LoginWindow</string>
	<key>RunAtLoad</key>
	<true/>
	<key>LaunchOnlyOnce</key>
	<true/>
	<key>StandardErrorPath</key>
	<string>/var/log/ds_finalize.log</string>
	<key>StandardOutPath</key>
	<string>/var/log/ds_finalize.log</string>
	<key>ProgramArguments</key>
	<array>
		<string>/etc/deploystudio/bin/ds_finalize.sh</string>
	</array>
</dict>
</plist>
