# Packages/Admin/DeployStudio Admin.app/Contents/Applications/DeployStudio Assistant.app/Contents/Resources/sysBuilder/10.8/LaunchDaemons/com.apple.DumpGPURestart.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.DumpGPURestart</string>
	<key>LaunchEvents</key>
	<dict>
		<key>com.apple.iokit.matching</key>
		<dict>
			<key>IOGPURestartReport</key>
			<dict>
				<key>IOProviderClass</key>
				<string>IOGPURestartReport</string>
			</dict>
                        <key>IOAccelGPURestartReport</key>
                        <dict>
                                <key>IOProviderClass</key>
                                <string>IOAccelGPURestartReport</string>
                        </dict>
		</dict>
	</dict>
	<key>ProgramArguments</key>
	<array>
		<string>/System/Library/Frameworks/OpenGL.framework/Versions/A/Libraries/DumpGPURestart</string>
	</array>
	<key>KeepAlive</key>
	<false/>
</dict>
</plist>
