# Packages/Admin/DeployStudio Admin.app/Contents/Applications/DeployStudio Assistant.app/Contents/Resources/sysBuilder/10.8/LaunchDaemons/com.apple.smb.preferences.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
        <!-- Run synchronize-preferences when the SMB server preferences are modified.  -->
        <key>Label</key>
        <string>com.apple.smb.preferences</string>
        <key>ProgramArguments</key>
        <array>
                <string>/usr/libexec/samba/smb-sync-preferences</string>
        </array>
        <key>EnableTransactions</key>
        <true/>
        <key>WatchPaths</key>
        <array>
            <string>/Library/Preferences/SystemConfiguration/com.apple.smb.server.plist</string>
        </array>
        <key>LaunchEvents</key>
        <dict>
            <!-- UserEventAgent ignores us unless we have a non-empty dictionary here. -->
            <key>com.apple.smb.preferences</key>
            <dict>
                <key>dummy</key>
                <true/>
            </dict>
        </dict>
        <key>TimeOut</key>
        <integer>1</integer>
</dict>
</plist>
