# Packages/Admin/DeployStudio Admin.app/Contents/Applications/DeployStudio Assistant.app/Contents/Resources/sysBuilder/10.8/LaunchDaemons/com.apple.NetworkDiagnostics.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple Computer//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Label</key>
	<string>com.apple.NetworkDiagnostics</string>
	<key>MachServices</key>
	<dict>
		<key>com.apple.NetworkDiagnostic.agent</key>
		<dict>
			<key>ResetAtClose</key>
			<true/>
		</dict>
	</dict>
	<key>Program</key>
	<string>/System/Library/CoreServices/Network Diagnostics.app/Contents/MacOS/Network Diagnostics</string>
	<key>ThrottleInterval</key>
	<integer>1</integer>
</dict>
</plist>
